# templates/systemd/opencode-serve.service
[Unit]
Description=OpenCode local server for Research-Code-Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/baseline-project
EnvironmentFile=/etc/research-code-agent/feishu-bridge.env
ExecStart=/usr/bin/env opencode serve --hostname 127.0.0.1 --port 4096
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
