# templates/systemd/rca-feishu-opencode-bridge.service
[Unit]
Description=Research-Code-Agent Feishu to OpenCode bridge
Wants=network-online.target
After=network-online.target opencode-serve.service
Requires=opencode-serve.service

[Service]
Type=simple
# Replace /opt/baseline-project and env path before installing.
WorkingDirectory=/opt/baseline-project
EnvironmentFile=/etc/research-code-agent/feishu-bridge.env
ExecStart=/usr/bin/env python3 /opt/baseline-project/tools/feishu_opencode_bridge.py
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
